# distro-files/layout/etc/systemd/system/prolinuxd.service
[Unit]
Description=Sineware ProLinux Daemon

[Service]
Type=simple
ExecStart=/opt/nodejs/bin/node /opt/prolinuxd/src/index.js
StandardOutput=kmsg+console
StandardError=kmsg+console
Environment=CONFIG_FILE=/sineware/data/prolinux.toml

RestartSec=10
Restart=always

[Install]
WantedBy=basic.target