# ocs2-prolinuxd/distro-files/prolinuxd.initd
#!/sbin/openrc-run
# ProLinuxD OpenRC Init file

command="node"
command_args="/opt/prolinuxd/index.js"
supervisor="supervise-daemon"

#command_background=true
#pidfile="/tmp/prolinuxd.pid"

output_log="/var/log/prolinuxd.log"
error_log="/var/log/prolinuxd.err"

depend() {
	want net
	before tinydm
}
